# yocto-next/net-rebooter/net-rebooter.service
[Unit]
Description=Check network connectivity check and reboot if needed

[Service]
Type=oneshot
ExecStart=/usr/local/bin/net-rebooter
